# crontab
# Extra security scanner
* * * * * /usr/bin/python3 /usr/local/bin/cyber-security-logger.py